# scripts/observability/com.verity.fairvalue-us-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  US fair-value(RIM V/P) 관측 로컬 weekly cron. 매주 토 08:00 KST (US 금요일 종가 후).
  Sharadar duckdb = 로컬 자산(repo 밖) → launchd 로컬 실행 (CI 아님, KR fair_value_gap 은 클라우드).
  산출 jsonl = own-use → 빌더가 .git-private(비공개 repo)에 self-commit+push (공개 main 차단).

  설치:
    cp scripts/observability/com.verity.fairvalue-us-weekly.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.verity.fairvalue-us-weekly.plist
  해제:
    launchctl unload ~/Library/LaunchAgents/com.verity.fairvalue-us-weekly.plist
  즉시 1회 실행(테스트):
    launchctl start com.verity.fairvalue-us-weekly

  StartCalendarInterval = 지정 시각 도래 시 실행. Mac off 면 누락(장기 trail 허용).
  로그: /tmp/fairvalue_us_weekly.out / .err. SEP 대용량 스캔으로 run ~3-5분.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.verity.fairvalue-us-weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/macbookpro/Desktop/배리티 터미널/scripts/observability/fair_value_gap_us.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/macbookpro/Desktop/배리티 터미널</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>6</integer>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/fairvalue_us_weekly.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/fairvalue_us_weekly.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
